# chrome-kiosk.service
[Unit]
Description=Chrome Kiosk Monitor
After=graphical.target

[Service]
ExecStart=/usr/local/bin/chrome-kiosk-monitor.sh
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
